# Modules/Other-clients/Client-UPS/configuration/Client-UPS/crontab
# run python script every 5mins
*/5 * * * * root /usr/local/bin/python /app/readUps.py > /app/logs/cron 2>&1

