# Client-N2_scraper/configuration/tank-scraper/crontab
# run python script every 15mins
*/15 * * * * root /usr/local/bin/python /app/readN2Tanks.py 2&>1
