# Client-N2_scraper/configuration/tank-scraper/crontab
# run python script every 10mins
*/10 * * * * root /usr/local/bin/python /app/readN2Tanks.py 2>&1
