# Other/MICS/MICS_Service.service
[Unit]
Description=Bruker MICS Monitoring

[Service]
ExecStart=/opt/Bruker/mics/runmicsctrl

[Install]
WantedBy=multi-user.target